# scripts/wsl/systemd/savh-worker-generate.service
[Unit]
Description=SAVH Print App - Worker (Generate)
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/savh_print_app
EnvironmentFile=/path/to/savh_print_app/.env
ExecStart=/usr/bin/env bash -lc 'poetry run python -m create_prints_server.worker.generate_worker'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

