# scripts/wsl/systemd/savh-worker-print.service
[Unit]
Description=SAVH Print App - Worker (Print)
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/savh_print_app
EnvironmentFile=/path/to/savh_print_app/.env
ExecStart=/usr/bin/env bash -lc 'poetry run python -m print_server.worker.print_worker'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target

